# server/crontab
1-59/2 * * * * ./~/hawk-track/server/execute.sh >/dev/null 2>&1
